# infra/systemd/tata-tcfzyq-api.service
[Unit]
Description=TATA tcfzyq API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/tata-h5/tcfzyq/current
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=8788
Environment=FAT_DB_PATH=/srv/tata-h5/tcfzyq/data/tata.sqlite
EnvironmentFile=-/etc/tata-tcfzyq-api.env
ExecStart=/usr/bin/node /srv/tata-h5/tcfzyq/current/apps/api/src/server.mjs
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/srv/tata-h5/tcfzyq/data

[Install]
WantedBy=multi-user.target
